# openvpn-ordinary-project-fixed/systemd/ovpnmgr-bot.service
[Unit]
Description=OpenVPN Manager Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ovpnmgr
Environment=PYTHONPATH=/opt/ovpnmgr/src
ExecStart=/opt/ovpnmgr/venv/bin/python -m ovpnmgr.bot
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
